# systemd/maschinenraum-whisper.service
[Unit]
Description=maschinenraum — vLLM Whisper STT server
After=network-online.target maschinenraum-core.service
Wants=network-online.target

[Service]
Type=simple
User=__MR_USER__
WorkingDirectory=__MR_HOME__
EnvironmentFile=__MR_HOME__/.env
ExecStart=__MR_HOME__/run.sh whisper
Restart=always
RestartSec=10
TimeoutStartSec=600
StandardOutput=append:__MR_HOME__/data/logs/whisper.log
StandardError=append:__MR_HOME__/data/logs/whisper.log

[Install]
WantedBy=multi-user.target
