# systemd/maschinenraum-tunnel.service
[Unit]
Description=maschinenraum — Cloudflare tunnel (remote access)
After=network-online.target maschinenraum-edge.service maschinenraum-core.service
Wants=network-online.target

[Service]
Type=simple
User=__MR_USER__
WorkingDirectory=__MR_HOME__
EnvironmentFile=__MR_HOME__/.env
ExecStart=__MR_HOME__/lib/cloudflare.sh run
Restart=always
RestartSec=10
StandardOutput=append:__MR_HOME__/data/logs/tunnel.log
StandardError=append:__MR_HOME__/data/logs/tunnel.log

[Install]
WantedBy=multi-user.target
