# systemd/maschinenraum-edge.service
[Unit]
Description=maschinenraum — Caddy router (one endpoint for chat + audio)
After=network-online.target maschinenraum-core.service
Wants=network-online.target

[Service]
Type=simple
User=__MR_USER__
WorkingDirectory=__MR_HOME__
ExecStart=/usr/bin/caddy run --config __MR_HOME__/Caddyfile --adapter caddyfile
Restart=always
RestartSec=5
StandardOutput=append:__MR_HOME__/data/logs/edge.log
StandardError=append:__MR_HOME__/data/logs/edge.log

[Install]
WantedBy=multi-user.target
