# systemd/maschinenraum-digest.service
[Unit]
Description=maschinenraum — weekly usage digest to Discord

[Service]
Type=oneshot
User=__MR_USER__
WorkingDirectory=__MR_HOME__
EnvironmentFile=__MR_HOME__/.env
ExecStart=__MR_HOME__/lib/metrics_digest.sh digest
